# scripts/launch/com.jace.crypto.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- LaunchAgent watchdog that checks the trading pipeline every 15 minutes -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jace.crypto.watchdog</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jacevanlierop/crypto_trading_bot/scripts/cron_watchdog.sh</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>900</integer>
    <key>StandardOutPath</key>
    <string>/Users/jacevanlierop/crypto_trading_bot/logs/watchdog.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jacevanlierop/crypto_trading_bot/logs/watchdog.launchd.err</string>
</dict>
</plist>
