# scripts/launch/com.jace.crypto.pipeline.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- LaunchAgent responsible for running the nightly trading pipeline on macOS -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jace.crypto.pipeline</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jacevanlierop/crypto_trading_bot/scripts/nightly_pipeline.sh</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/jacevanlierop/crypto_trading_bot</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>5</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/jacevanlierop/crypto_trading_bot/logs/pipeline.launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/jacevanlierop/crypto_trading_bot/logs/pipeline.launchd.err</string>
</dict>
</plist>
